# supervisor_minidlna.conf
[program:minidlna]
user = root 
command = /usr/sbin/minidlnad -S -R 
